# files/etc/cron.d/easypress_system_checks.master
10 */4 * * *  root /usr/local/sbin/check_size_of_php_error_logs.sh
