# files/etc/cron.d/eqpress_backups.slave
26 3,7,11,15,19,23      * * *   root    /usr/local/sbin/backup_from_master.sh > /dev/null
