# files/etc/cron.hourly/ntpdate
#!/bin/sh
/usr/sbin/ntpdate -u -b -s 132.246.168.148 142.3.100.15 132.246.168.164
