# files/etc/cron.hourly/ntpdate
#!/bin/sh
/usr/sbin/ntpdate -u -b -s 0.pool.ntp.org 1.pool.ntp.org 2.pool.ntp.org
