# roles/timed-jobs/files/etc/cron.d/eqpress_system_checks.master
10 * * * *    root /usr/local/sbin/check_for_large_http_log_files.sh
20 */4 * * *  root /usr/local/sbin/check_for_disk_usage.sh
10 */4 * * *  root /usr/local/sbin/check_size_of_php_error_logs.sh
