# roles/timed-jobs/files/etc/cron.d/eqpress_system_checks.slave
1 */6 * * *  root /usr/local/sbin/check_for_large_files.sh
2 */4 * * *  root /usr/local/sbin/check_for_disk_usage.sh
*/5 * * * *  root /usr/local/sbin/check_mysql_replication_status.sh
3 */4 * * *  root /usr/sbin/nginx -s reload
