# roles/mariadb-local-service/templates/etc/systemd/system/mariadb.service.d/execstartpost.conf.j2
# {{ ansible_managed }}
[Service]
ExecStartPost={{ mariadb_execstartpost_override|default() }}
