# roles/mariadb-local-service/templates/etc/systemd/system/mariadb.service.d/limits.conf.j2
# {{ ansible_managed }}
[Service]
LimitNOFILE={{ mariadb_systemd_override_nofile }}
LimitMEMLOCK={{ mariadb_systemd_override_memlock }}
