# roles/eqpress-cronjobs/files/etc/cron.d/eqpress_backups.master
#13 4    * * *   root    /usr/local/sbin/wp_db_backup.sh
13 4    * * *   root    /usr/local/sbin/ep_site_snapshot.sh
